# exchange-analytics-scheduler.service
[Unit]
Description=Exchange Analytics Data Scheduler
Documentation=https://github.com/m37335/exchangeApp
After=network.target
Wants=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/app
Environment=PATH=/usr/local/bin:/usr/bin:/bin
ExecStart=/usr/bin/python3 /app/data_scheduler.py
ExecReload=/bin/kill -HUP $MAINPID
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ReadWritePaths=/app/logs

# Resource limits
LimitNOFILE=1024
MemoryMax=512M

[Install]
WantedBy=multi-user.target
